# config/autostart/mouseless.desktop
[Desktop Entry]
Type=Application
Name=mouseless

Path=/var/home/sofusa/.local/bin
Exec=mouseless
Terminal=false

